# crontab
*/5 * * * *    /usr/local/bin/python /usr/local/bin/lnetatmo_influx.py >> /var/log/cron.log 2>&1
